# docker/systemd/grammashop-backup.service
[Unit]
Description=grammashop: nightly Postgres backup to S3
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/shmonke/grammashop
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml exec -T api node apps/api/dist/db/backup.js
